# deploy/roles/Stouts.wsgi/templates/upstart.conf.j2
#!upstart

# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

description "Start uwsgi for {{wsgi_app_name}}"

start on (local-filesystems and runlevel [2345])
stop on runlevel [!2345]

# expect fork
respawn

limit nofile 65536 65536
console log
chdir {{wsgi_src_dir}}
setuid {{wsgi_user}}
setgid {{wsgi_group}}

{% if wsgi_env_dir %}
exec {{wsgi_env_dir}}/bin/uwsgi --ini {{wsgi_uwsgi_ini}}
{% else %}
exec /usr/local/bin/uwsgi --ini {{wsgi_uwsgi_ini}}
{% endif %}

