# deploy/roles/Stouts.celery/templates/upstart.conf.j2
#!upstart

# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

description "Start celery workers for {{celery_app_name}}"

start on (local-filesystems and runlevel [2345])
stop on runlevel [!2345]

{% if item.user|default(celery_user) == 'root' %}
env C_FORCE_ROOT=true
env HOME=/root
{% endif %}

env USER="{{ item.user|default(celery_user) }}"
env GROUP="{{ item.group|default(celery_group) }}"

env CELERY_NAME=""
env CELERY_APP="{{ item.app_module|default(celery_app_module) }}"
env CELERY_CHDIR="{{ item.app_dir|default(celery_app_dir) }}"
env CELERY_LOG_FILE="{{ item.log_dir|default(celery_log_dir) }}/{{item.action|default('worker')}}{{item.queue|default('') and ('-' + item.queue) or ''}}.log"
env CELERY_LOG_LEVEL="{{ item.log_level|default(celery_log_level) }}"
env CELERY_PID_FILE="{{ item.run_dir|default(celery_run_dir) }}/{{item.action|default('worker')}}{{item.queue|default('') and ('-' + item.queue) or ''}}.pid"

env CELERY_QUEUE_ARG="{{ '-Q %s' % item.queue if item.queue|default(None) else '' }}"
env CELERY_CONCURRENCY_ARG="{{ '-c %s' % item.concurrency|default(celery_concurrency) if item.action == 'worker' else '' }}"

console log
kill timeout 20
limit nofile 65536 65536
respawn

exec {{celery_bin}} {{item.action}} $CELERY_QUEUE_ARG $CELERY_CONCURRENCY_ARG \
                                    --pidfile="$CELERY_PID_FILE" \
                                    --logfile="$CELERY_LOG_FILE" \
                                    --loglevel="$CELERY_LOG_LEVEL" \
                                    --app="$CELERY_APP" \
                                    --workdir="$CELERY_CHDIR" \
                                    --uid="$USER" \
                                    --gid="$GROUP"
